# srcs/wordpress/srcs/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
nodaemon=true
user=root

[program:nginx]
command=sh -c "kill $(pidof nginx) || nginx -g 'daemon off;' && kill $(pidof supervisord)"
;command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/nginx.log
autostart=true
autorestart=true
stopsignal=KILL
numprocs=1

[program:php-fpm]
command=sh -c "kill $(pidof php-fpm7) || /usr/sbin/php-fpm7 -R -F && kill $(pidof supervisord)"
;command=/usr/sbin/php-fpm7 -R -F
autostart=true
autorestart=true
