# srcs/grafana/srcs/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
nodaemon=true

[program:grafana]
command=sh -c "./grafana-server && kill -s SIGTERM $(cat supervisord.pid)"
autostart=true
autorestart=true

[program:telegraf]
command=sh -c "/usr/bin/telegraf --config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"
autostart=true
autorestart=true
