# srcs/nginx/srcs/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
nodaemon=true
user=root

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true